# supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/init.sh
autorestart=false
priority=1

[program:postgres]
command=/usr/lib/postgresql/13/bin/postgres -D /postgresql
priority=2
autorestart=true
user=postgres

[program:mongo]
command=mongod --bind_ip 0.0.0.0 --replSet rs0
autorestart=true

[program:redis]
command=redis-server
autorestart=true

[program:clickhouse]
command=clickhouse-server --config /etc/clickhouse-server/config.xml
autorestart=true
user=clickhouse
stdout_logfile=/var/log/supervisor/clickhouse-stdout.log
stderr_logfile=/var/log/supervisor/clickhouse-stderr.log

[program:zookeeper]
command=/opt/kafka/bin/zookeeper-server-start.sh /opt/kafka/config/zookeeper.properties
autorestart=true

[program:kafka]
command=/opt/kafka/bin/kafka-server-start.sh /opt/kafka/config/server.properties
autorestart=true

[program:dbeaver]
command=/opt/cloudbeaver/run-server.sh 
directory=/opt/cloudbeaver
autorestart=true
environment=
    CB_SERVER_PORT="8978",
    CB_SERVER_PUBLIC_URL="0.0.0.0:8978",
    JAVA_HOME="/usr/lib/jvm/java-17-openjdk-amd64"

[program:spark]
command=/opt/spark/sbin/start-master.sh
autorestart=false
environment=
    SPARK_MASTER_HOST=localhost,
    SPARK_IDENT_STRING=master,
    SPARK_MASTER_PORT=7077,
    SPARK_MASTER_WEBUI_PORT=8081,
    SPARK_WORKER_PORT=8881

[program:airflow]
command=airflow webserver
autorestart=true
dependent_startup_wait_for=postgres:running
dependent_startup=true
environment=
    AIRFLOW_HOME="/root/airflow",
    AIRFLOW__DATABASE__SQL_ALCHEMY_CONN="postgresql+psycopg2://postgres:postgres@localhost:5432/airflow",
    AIRFLOW__CORE__EXECUTOR="LocalExecutor"

[program:airflow-scheduler]
command=airflow scheduler
autorestart=true
environment=
    AIRFLOW_HOME="/root/airflow",
    AIRFLOW__DATABASE__SQL_ALCHEMY_CONN="postgresql+psycopg2://postgres:postgres@localhost:5432/airflow",
    AIRFLOW__CORE__EXECUTOR="LocalExecutor"

[program:superset]
command=superset run -h 0.0.0.0 -p 8089
autorestart=true
environment=
    FLASK_APP="superset.app:create_app()",
    SUPERSET_SECRET_KEY="$(openssl rand -base64 42)",
    SUPERSET_CONFIG_PATH=/superset_config.py
stdout_logfile=/var/log/supervisor/superset-stdout.log
stderr_logfile=/var/log/supervisor/superset-stderr.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log
user=root